# deploy/clientst0r-breach-scan.timer
[Unit]
Description=Run Client St0r password breach scanning daily
Requires=clientst0r-breach-scan.service

[Timer]
# Run daily at 02:00 local time
OnCalendar=*-*-* 02:00:00

# Run on boot if missed (host was off at 02:00)
Persistent=true

# Random delay up to 30 min to avoid HIBP API peak times
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
