# deploy/clientst0r-rmm-sync.timer
[Unit]
Description=Client St0r RMM Sync Timer
Requires=clientst0r-rmm-sync.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
