# deploy/clientst0r-breach-scan.timer
[Unit]
Description=Run Client St0r password breach scanning daily
Requires=clientst0r-breach-scan.service

[Timer]
# Run daily at 2 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00

# Run on boot if missed
Persistent=true

# Random delay to avoid peak times
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
